# wedding-api.service
[Unit]
Description=Wedding Site API
After=network.target

[Service]
Type=simple
User=www-data
WorkingDirectory=/var/www/wedding-site
ExecStart=/usr/bin/node server.js
Restart=on-failure
Environment=NODE_ENV=production
Environment=PORT=3000

[Install]
WantedBy=multi-user.target 